# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0

# --- Asterisk: the PBX that terminates the Loxone SIP call --------------------
[program:asterisk]
command=/usr/sbin/asterisk -f -vvv -C /etc/asterisk/asterisk.conf
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# --- Bridge: connects Protect <-> Asterisk; also spawns go2rtc ----------------
# Waits for ARI (8088) before starting, so it doesn't race Asterisk on boot.
[program:bridge]
command=/bin/bash -c 'until curl -sf -u "$ARI_USER:$ARI_PASS" http://127.0.0.1:8088/ari/asterisk/info >/dev/null 2>&1; do echo "[bridge] waiting for ARI..."; sleep 2; done; echo "[bridge] ARI up, starting"; exec node /app/index.js'
directory=/app
autostart=true
autorestart=true
priority=20
startsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# --- Video transcode: continuous ffmpeg -> /dev/shm/cam.jpg (always warm) -----
# Reads go2rtc's local RTSP (spawned by the bridge), so starts after it; its own
# retry loop + autorestart cover go2rtc not being up yet.
[program:mjpeg-video]
command=/usr/local/bin/mjpeg-video.sh
autostart=true
autorestart=true
priority=30
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# --- mjpg-streamer: serves the JPEG at :MJPEG_PORT/?action=stream (+ snapshot),
# with optional basic auth. Wrapper reads MJPEG_PORT / MJPEG_AUTH_* from env.
[program:mjpg-streamer]
command=/usr/local/bin/mjpg-serve.sh
autostart=true
autorestart=true
priority=30
startsecs=5
startretries=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# --- camera-rotation-watchdog: keeps the doorbell LANDSCAPE-UPRIGHT -----------
# The doorbell firmware forces a 90deg portrait rotation Protect can't disable.
# This SSHes to the camera every ROTATION_CHECK_INTERVAL s and re-applies an
# LD_PRELOAD hook (hallwayMode=0 -> landscape) + ISP 180 flip whenever it drifts
# (e.g. after a camera reboot). No-op unless ROTATION_FIX_ENABLED=1. Idles (never
# exits) when disabled/unconfigured, so supervisor doesn't flap it.
[program:rotation-watchdog]
command=/app/patches/camera-rotation-watchdog.sh
autostart=true
autorestart=true
priority=40
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
